# laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/#user#/#project-name#/artisan queue:listen --tries=3           
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=#user#
numprocs=1
redirect_stderr=true
stdout_logfile=/home/#user#/#project-name#/storage/logs/worker.log
stopwaitsecs=30

#Path can be retrieved from Root Path#
